# scripts/resy-monitor.service
[Unit]
Description=Resy Sniper - Monitor Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/resy-sniper
EnvironmentFile=/home/ec2-user/resy-sniper/.env
ExecStart=/home/ec2-user/resy-sniper/.venv/bin/python -m monitor.main
Restart=on-failure
RestartSec=10
Environment=PYTHONPATH=/home/ec2-user/resy-sniper

[Install]
WantedBy=multi-user.target
